# sparse/etc/systemd/system/sensorfwd.service.d/99-wait-droid-hal.conf
[Unit]
After=droid-hal.service
Wants=droid-hal.service

[Service]
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
