# sparse/usr/lib/systemd/system/nemo-devicelock.socket.d/99-runtime-dir.conf
[Socket]
RuntimeDirectory=nemo-devicelock
RuntimeDirectoryMode=0755
