# sparse/usr/lib/systemd/system/vendor-bt_firmware.mount
[Unit]
Description=Droid mount for /vendor/bt_firmware
After=vendor.mount
Requires=vendor.mount
Before=local-fs.target

[Mount]
What=/dev/sde24
Where=/vendor/bt_firmware
Type=vfat
Options=ro,shortname=lower,uid=1002,gid=3002,dmask=227,fmask=337,context=u:object_r:bt_firmware_file:s0
TimeoutSec=5

[Install]
WantedBy=local-fs.target
